# rootoverlay/etc/systemd/system/multi-user.target.wants/man-db.service
[Unit]
Description=Update man-db cach

[Service]
Type=disabled
